# roles/apt_mirror/files/apt-mirror.service
[Unit]
Description=Apt Mirror Repository Synchronization
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/apt-mirror
User=root

